# akita/etc/cron.d/akita
SHELL=/bin/bash
0 3,6,9 * * * akita source /home/akita/.bashrc && cd /var/www/akita/current && bin/rake salesforce:import_users >> /var/www/akita/current/log/salesforce_import.log
